# home/etc/profile.d/java.sh
export JAVA_HOME=/opt/jse7
export M2_HOME=/opt/maven3
export GROOVY_HOME=/opt/groovy2
export ANT_HOME=/opt/ant1
export GRADDLE_HOME=/opt/graddle
export PATH=${JAVA_HOME}/bin:${M2_HOME}/bin:${GROOVY_HOME}/bin:${ANT_HOME}/bin:${GRADDLE_HOME}/bin:${PATH}
